# ansible/files/signage-cloud-sync.timer
[Unit]
Description=Run cloud synchronization regularly

[Timer]
OnActiveSec=1min
OnUnitInactiveSec=30min

[Install]
WantedBy=timers.target
