# ansible/files/signage-video.timer
[Unit]
Description=Run video regularly

[Timer]
OnActiveSec=1min
OnUnitInactiveSec=30sec
AccuracySec=1sec

[Install]
WantedBy=timers.target
