# ansible/files/signage-video-loop-restart.path
[Path]
PathChanged=/opt/svpb-signage/content/video.mp4

[Install] 
WantedBy=xsession.target
