# ansible/files/xsession.target
[Unit]
Description=Pseudo unit to know when xsession is running
BindsTo=graphical-session.target
